# rootfs/.init/services/test.conf
[program:test]
command = echo "this is a test"
startsecs = 0
autorestart = false
startretries = 1